# ops/supervisor.smartpet-worker.conf
[program:smartpet-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/smartpet/backend/artisan queue:work --sleep=1 --tries=3 --timeout=120
autostart=true
autorestart=true
user=deploy
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/smartpet-worker.log
stopwaitsecs=3600

